# home/dot_config/systemd/user/redshift.service.d/override.conf
[Unit]
After=graphical-session.target
Requires=graphical-session.target
